# systemd/ipfs.service.d/autotls.conf
[Service]
Environment=GOLOG_LOG_LEVEL=error,autotls=info
